# zika-server.service
[Unit]
After=network.target

[Service]
WorkingDirectory=<...>
ExecStart=<...>/zika-server
Restart=always
